# distrib/debian-3.1/preinst
#!/bin/sh
invoke-rc.d steam stop

if [ "$action" = upgrade ] || [ "$action" = install ]
then
  if ! test -f /usr/share/debconf/confmodule
  then
    cat <<EOF

WARNING: debconf is not available

Please read /usr/share/doc/steam/README.Debian for more information how to
setup steam by hand.

EOF
    exit 0
  else
    . /usr/share/debconf/confmodule
    db_version 2.0
  fi
fi

# dh_installdeb will replace this with shell code automatically
# generated by Aother debhelper scripts.
#DEBHELPER#

exit 0

